# python/backend/docker/supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; supervisord.conf — manages Gunicorn workers inside the container
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
user=appuser

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Gunicorn ─────────────────────────────────────────────────────────────────
[program:gunicorn]
command=gunicorn service_platform.wsgi:application
    --bind 0.0.0.0:8000
    --workers %(ENV_GUNICORN_WORKERS)s
    --worker-class gthread
    --threads 2
    --timeout 120
    --keep-alive 5
    --max-requests 1000
    --max-requests-jitter 50
    --log-level info
    --access-logfile -
    --error-logfile -
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DJANGO_SETTINGS_MODULE="service_platform.settings"
