# golang/backend/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
priority=10                     ; start nginx first

; ── Go API ────────────────────────────────────────────────────────────────────
[program:goapi]
command=/app/server
autostart=true
autorestart=true
startretries=5
startsecs=2                     ; must stay up for 2s to be considered started
stdout_logfile=/var/log/supervisor/goapi.log
stderr_logfile=/var/log/supervisor/goapi-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=20                     ; start after nginx

; ── supervisorctl web interface (optional, remove if not needed) ───────────────
[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001
